# template/supervisord.d/sixpack-web.ini
[program:sixpack-web]
command = bin/sixpack-web
autorestart = true
stdout_logfile = /var/log/supervisor/sixpack-web.log
redirect_stderr = true
environment = PYTHONPATH=".",SIXPACK_CONFIG="./config.yml"
directory = /home/sixpack/sixpack
