# template/supervisord.d/sixpack.ini
[program:sixpack]
command = bin/sixpack
autorestart = true
stdout_logfile = /var/log/supervisor/sixpack.log
redirect_stderr = true
environment = PYTHONPATH=".",SIXPACK_CONFIG="./config.yml"
directory = /home/sixpack/sixpack
